# supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:nginx]
command = /init.sh nginx
autostart = True

[program:python]
command = python3 /geo-python/main.py
autostart = True